# scripts/podman/propel.container.in
# Propel gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into propel's ~/.config/containers/systemd/
# {{PROPEL_HOME}} is replaced at install time.

[Unit]
Description=Propel gateway (rootless Podman)

[Container]
Image=propel:local
ContainerName=propel
UserNS=keep-id
Volume={{PROPEL_HOME}}/.propel:/home/node/.propel
EnvironmentFile={{PROPEL_HOME}}/.propel/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
